# systemd/bonefire-tunnel.service
[Unit]
Description=Bonefire ngrok Tunnel
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/bonefire
ExecStart=/opt/bonefire/bonefire_env/bin/python /opt/bonefire/bonefire_tunnel.py
Restart=always

[Install]
WantedBy=multi-user.target
