# systemd/bonfire-tunnel.service
[Unit]
Description=Bonfire ngrok Tunnel
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/bonfire
ExecStart=/usr/bin/env python3 /opt/bonfire/bonfire_tunnel.py
Restart=always

[Install]
WantedBy=multi-user.target
